# frontend/deploy/systemd/auusite-frontend.service
[Unit]
Description=AAU Frontend (Next.js)
After=network.target

[Service]
Type=simple
User=frappe
WorkingDirectory=/home/frappe/frappe-bench/apps/aau_university/AAU-16.24
Environment=NODE_ENV=production
Environment=PATH=/home/frappe/.nvm/versions/node/v20.19.0/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/home/frappe/.nvm/versions/node/v20.19.0/bin/npm run start:hosted
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
